# rupiyame-fullstack.service
[Unit]
Description=RupiyaMe CRM Full Stack Service (Backend + Frontend)
After=network.target mongodb.service
Wants=network.target mongodb.service

[Service]
Type=forking
User=ubuntu
Group=ubuntu
WorkingDirectory=/www/wwwroot/rupiyamaker/RupiyaMe
ExecStart=/www/wwwroot/rupiyamaker/RupiyaMe/start-full-crm.sh start
ExecStop=/www/wwwroot/rupiyamaker/RupiyaMe/start-full-crm.sh stop
ExecReload=/www/wwwroot/rupiyamaker/RupiyaMe/start-full-crm.sh restart
Restart=always
RestartSec=10
TimeoutStartSec=120

# Security enhancements
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target