# rupiyame-updated.service
[Unit]
Description=RupiyaMe CRM Backend Service
After=network.target mongodb.service
Wants=network.target mongodb.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/www/wwwroot/rupiyamaker/RupiyaMe/backend
ExecStart=/bin/bash -c 'source /www/wwwroot/rupiyamaker/RupiyaMe/backend/venv/bin/activate && ENVIRONMENT=development python3 -m app'
Restart=always
RestartSec=10
Environment=PYTHONPATH=/www/wwwroot/rupiyamaker/RupiyaMe/backend
Environment=ENVIRONMENT=development

# Increase timeout for slower systems or heavy initialization
TimeoutStartSec=90
TimeoutStopSec=30

# Security enhancements
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rupiyame

[Install]
WantedBy=multi-user.target